# deploy/bot-ii.service
[Unit]
Description=Telegram AI bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/go/bot-ii
EnvironmentFile=/var/go/bot-ii/.env
Environment=HTTPS_PROXY=http://127.0.0.1:2081
Environment=NO_PROXY=127.0.0.1,localhost
ExecStart=/var/go/bot-ii/bot-ii
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
